# sup_freeswitch.conf
[program:FreeSwitch]
command=/usr/local/freeswitch/bin/freeswitch -nc
autostart=true
autorestart=true